# deploy/sonoglyph-bridge.service
[Unit]
Description=Sonoglyph bridge (Hono + WebSocket + MCP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/sonoglyph/proxy
EnvironmentFile=/srv/sonoglyph/.env
Environment=NODE_ENV=production
Environment=PROXY_PORT=8787
Environment=PUBLIC_BASE_URL=https://sonoglyph.xyz
Environment=CORS_ORIGINS=https://sonoglyph.xyz,https://www.sonoglyph.xyz
ExecStart=/srv/sonoglyph/proxy/node_modules/.bin/tsx src/index.ts
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal
# Hardening — bridge doesn't need filesystem write outside /srv/sonoglyph
# and /tmp (via PrivateTmp). ProtectHome is `read-only` (not `true`) so the
# bridge can read /root/.hermes/.env + config.yaml when spawning ephemeral
# shared-agent Hermes processes, AND exec the operator's uv-managed Python
# interpreter at /root/.local/share/uv/python/<…>/bin/python3.11 (which is
# what /usr/local/bin/hermes's shebang chain resolves to). With
# ProtectHome=true, both paths returned ENOENT to the service and shared-
# agent spawn failed with exit code 127.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/srv/sonoglyph
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
